# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/app/nginx.log
stderr_logfile=/var/log/app/nginx_error.log

[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8081 --workers 1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/app/uvicorn.log
stderr_logfile=/var/log/app/uvicorn_error.log
environment=AICP_DATA_DIR="/data/results",AICP_PORT="8081"

[program:sshd]
command=/usr/sbin/sshd -D
autostart=%(ENV_SSH_ENABLED)s
autorestart=true
stdout_logfile=/var/log/app/sshd.log
stderr_logfile=/var/log/app/sshd_error.log
